# sv/suricata/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars
env

:${INTERFACE="eth0"}

exec 2>&1
/usr/bin/suricata -c /etc/suricata//suricata.yaml -i ${INTERFACE}
